# deb/debian/homebridge.service
[Unit]
Description=Homebridge
Wants=network-online.target
After=syslog.target network-online.target

[Service]
Type=simple
User=homebridge
WorkingDirectory=/var/lib/homebridge
ExecStart=/opt/homebridge/start.sh
Restart=always
RestartSec=3
KillMode=process

[Install]
Alias=homebridge
WantedBy=multi-user.target
